# docker/supervisord/mysql.conf
[program:mysqld]
command=/usr/bin/mysqld_safe
numprocs=1
autostart=true
autorestart=true
startretries=1
startsecs=1
stderr_logfile=/var/log/mysql/supervisor_stderr.log
stdout_logfile=/var/log/mysql/supervisor_stdout.log
